# supervisor.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/cli.log
stderr_logfile=/cli.err

[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/cli.log
stderr_logfile=/cli.err

[program:php-monitor]
command=php /extra/monitor/start.php
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/cli.log
stderr_logfile=/cli.err

[program:php-init]
command=php /extra/monitor/init.php
autostart=true
autorestart=false
startsecs=3
stdout_logfile=/cli.log
stderr_logfile=/cli.err