# supervisor/init.conf
[program:php-init]
command=php /extra/init.php
autostart=true
autorestart=false
startsecs=1
startretries=1
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0