# supervisor/monitor.conf
[program:php-monitor]
command=php /extra/monitor.php
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0